# templates/myapp.service.j2
[Unit]
Description=myapp
After=multi-user.target

[Service]
Type=simple
ExecStart=python3 /usr/share/myapp.py

[Install]
WantedBy=multi-user.target
